# package/czechlight-rauc/czechlight.target
[Unit]
Description=A fully functional CzechLight device
Requires=multi-user.target enable-hw-watchdog.service rauc-mark-good.service
