# package/czechlight-cfg-fs/sysrepo-persistent-cfg.service
[Unit]
Description=Persisting persistent sysrepo datastores to /cfg
After=netopeer2.service cfg.mount
Requires=netopeer2.service cfg.mount

[Service]
Type=simple
ExecStart=/bin/sh -c 'while true; do inotifywait -e CLOSE_WRITE /etc/sysrepo/data/*.startup && mkdir -p /cfg/sysrepo-startup/ && cp -a /etc/sysrepo/data/*.startup /cfg/sysrepo-startup/; done'
