# package/czechlight-cfg-fs/sysrepo-persistent-cfg.service
[Unit]
Description=Persisting persistent sysrepo datastores to /cfg
After=netopeer2.service cfg.mount
Requires=netopeer2.service cfg.mount

[Service]
Type=simple
UMask=0077
ExecStart=/bin/sh -c 'while true; do inotifywait -e CLOSE_WRITE /etc/sysrepo/data/*.startup && mkdir -p /cfg/sysrepo/ && sysrepocfg -d startup -f json -X > /cfg/sysrepo/startup.json; done'
Group=sysrepo
