# package/czechlight-cfg-fs/yang-startup.target
[Unit]
Description=All YANG models installed, all data needed, all migrations done
