# package/cla-sysrepo/cla-appliance.service.in
[Unit]
Description=CzechLight __MODEL__ driver
After=syslog.target network.target cfg-restore-sysrepo.service czechlight-migrate.service nacm-restore.service
Before=rauc-mark-good.service velia-hardware-g1.service velia-hardware-g2.service
PartOf=netopeer2.service
Requires=cfg-restore-sysrepo.service czechlight-migrate.service nacm-restore.service
StartLimitIntervalSec=0
ConditionKernelCommandLine=|czechlight=__MODEL__
ConditionKernelCommandLine=|czechlight=__MODEL__-g2

[Service]
Type=notify
ExecStart=/usr/bin/cla-sysrepod --io-log-level=5 --properties-log-level=5 --sr-bridge-log-level=5 --sysrepo-log-level=3 --appliance=__MODEL__
Group=sysrepo
PrivateTmp=yes
PrivateDevices=no
ProtectSystem=yes
ProtectHome=yes
Restart=on-abnormal
SystemCallArchitectures=native
Restart=always
WatchdogSec=15s
LogRateLimitIntervalSec=10
LogRateLimitBurst=30000
SyslogLevel=alert

[Install]
WantedBy=multi-user.target
RequiredBy=velia-hardware-g1.service velia-hardware-g2.service
