# package/czechlight-power/czechlight-power.service
[Unit]
Description=Dynamic loading of power supply drivers

[Service]
ExecStart=/usr/bin/czechlight-power.sh
